# shadowsocks/ss-subscribe-update.service
[Unit]
Description=Daily Update Subscribe Service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/ss-subscribe.sh update-subscribe
